# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
exec -g "daemon off;" /usr/sbin/nginx -c /etc/nginx/nginx.conf
